# deploy/systemd/lms-daily-booster.timer
[Unit]
Description=Create daily LMS Brain Booster reminders at 09:00 IST

[Timer]
OnCalendar=*-*-* 09:00:00 Asia/Kolkata
Persistent=true
AccuracySec=1min
Unit=lms-scheduled-job@daily-booster.service

[Install]
WantedBy=timers.target
